# systemd/user/display_information.service
[Unit]
Description=Display ip address and battery information on AtomS3
Wants=network-online.target
After=network-online.target

[Service]
User=%I
Group=%I
Environment="HOME=/home/%I"
EnvironmentFile=-/etc/opt/riberry/riberry.user
Environment="LOG_LEVEL=DEBUG"
Type=simple
ExecStart=/usr/local/bin/run_display_information.sh
Restart=always
RestartSec=10
SyslogIdentifier=%n

[Install]
WantedBy=multi-user.target
